# rootfs/etc/services.d/graylog2/run
#!/usr/bin/with-contenv sh

/opt/graylog/bin/graylogctl run